# systemd/sector-dashboard.service
# Sector Rotation Dashboard - systemd unit file
#
# Install:
#   sudo cp systemd/sector-dashboard.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now sector-dashboard
#
# This template targets the AHADPI5 deployment checkout. If you install it on a
# different host, adjust User=, Group=, WorkingDirectory=, and ExecStart=.

[Unit]
Description=Sector Rotation Dashboard (Streamlit)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ahad
Group=ahad
WorkingDirectory=/home/ahad/SECTOR_MOMENTUM_AND_ROTATION
Environment=STATE_FILE=/home/ahad/SECTOR_MOMENTUM_AND_ROTATION/data/state.json
Environment=STATE_TRANSITION_JOURNAL=/home/ahad/SECTOR_MOMENTUM_AND_ROTATION/data/state_transitions.jsonl
ExecStart=/home/ahad/SECTOR_MOMENTUM_AND_ROTATION/.venv/bin/streamlit run app.py \
    --server.address 127.0.0.1 \
    --server.port 8501 \
    --server.headless true \
    --server.fileWatcherType none \
    --browser.gatherUsageStats false
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Reasonable defaults
TimeoutStartSec=120
TimeoutStopSec=30
KillMode=mixed
KillSignal=SIGINT

[Install]
WantedBy=multi-user.target
