# systemd/user/sector-provider-flow-cache.service
# Sector Rotation Provider-Flow Cache Warmup - user systemd one-shot unit
#
# Install without sudo:
#   mkdir -p ~/.config/systemd/user
#   cp systemd/user/sector-provider-flow-cache.service ~/.config/systemd/user/
#   cp systemd/user/sector-provider-flow-cache.timer ~/.config/systemd/user/
#   systemctl --user daemon-reload
#   systemctl --user enable --now sector-provider-flow-cache.timer

[Unit]
Description=Sector Rotation Provider-Flow Cache Warmup
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=%h/SECTOR_MOMENTUM_AND_ROTATION
Environment=PYTHONUNBUFFERED=1
ExecStart=%h/SECTOR_MOMENTUM_AND_ROTATION/.venv/bin/python scripts/warm_provider_flow_cache.py --universe scored --timeout 20
